# laravel-worker2.conf
[program:laravel-worker2]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work --queue=second --sleep=3 --tries=2
autostart=true
autorestart=true
user=root
numprocs=4
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/supervisor_queue-work2.log

